# game-runner.service
[Unit]
Description=Run a Dragonshark game

[Service]
User=root
ExecStart=/bin/game-runner

[Install]
WantedBy=multi-user.target